# rootfs/etc/s6-overlay/s6-rc.d/svc-automation/run
#!/command/with-contenv bash
# shellcheck shell=bash
set -e

CDP_PORT="${CHROME_DEBUGGING_PORT:-9229}"
CDP_URL="http://localhost:${CDP_PORT}"

echo "[svc-automation] Waiting for Chrome debugging port..."

# Wait for Chrome CDP to be available
timeout=120
while [ $timeout -gt 0 ]; do
    if curl -sf "${CDP_URL}/json/version" >/dev/null 2>&1; then
        echo "[svc-automation] Chrome CDP is ready!"
        break
    fi
    sleep 1
    timeout=$((timeout - 1))
    if [ $((timeout % 10)) -eq 0 ]; then
        echo "[svc-automation] Still waiting for Chrome CDP... ($timeout seconds left)"
    fi
done

if [ $timeout -eq 0 ]; then
    echo "[svc-automation] ERROR: Chrome debugging port not available after 120s"
    exit 1
fi

echo "[svc-automation] Starting Python automation bot..."
cd /app/automation

# Run Python with unbuffered output (-u) and force stdout/stderr to be line-buffered
# Using exec so s6 properly supervises the Python process
exec stdbuf -oL -eL python3 -u bot.py 2>&1
