# rootfs/etc/s6-overlay/s6-rc.d/init-pulseaudio-setup/up
/etc/s6-overlay/scripts/pulseaudio-setup
