# rootfs/etc/s6-overlay/s6-rc.d/svc-pulseaudio-custom/run
#!/command/with-contenv bash
# shellcheck shell=bash
# Custom PulseAudio service — closely matches original Nomad audio-hub task
set -e

echo "[svc-pulseaudio-custom] Starting PulseAudio..."

# Create socket directory (matching original Nomad config)
mkdir -p /alloc/pulse
rm -f /alloc/pulse/socket

# Create a runtime directory PulseAudio can actually use
mkdir -p /var/run/pulse
chmod 755 /var/run/pulse

# Set HOME to a writable location (the original Nomad container
# ran in debian:bookworm-slim where /root was writable; in this
# linuxserver image it isn't)
export HOME=/var/run/pulse
export XDG_RUNTIME_DIR=/var/run/pulse

# Run PulseAudio in user mode (matching original Nomad audio-hub task)
# The original used: pulseaudio --start --daemonize=no --exit-idle-time=-1 --verbose
# We skip --start since s6 manages the lifecycle
exec pulseaudio \
    --daemonize=no \
    --exit-idle-time=-1 \
    --verbose \
    --log-target=stderr \
    --load="module-native-protocol-unix socket=/alloc/pulse/socket auth-anonymous=1" \
    --load="module-null-sink sink_name=VirtualMic sink_properties=device.description=Input_To_ChatGPT" \
    --load="module-null-sink sink_name=VirtualSpeaker sink_properties=device.description=Output_From_ChatGPT" \
    --load="module-remap-source master=VirtualMic.monitor source_name=RealMic source_properties=device.description=ChatGPT_Microphone" \
    2>&1
