# rootfs/etc/s6-overlay/s6-rc.d/svc-baresip/finish
#!/command/with-contenv bash
# shellcheck shell=bash

# Clean up call_active file on shutdown
rm -f /alloc/call_active
echo "[svc-baresip] Service stopped, cleaned up call_active file"
