# rootfs/etc/s6-overlay/s6-rc.d/svc-baresip/run
#!/command/with-contenv bash
# shellcheck shell=bash
set -e

# PULSE_SERVER is set in Dockerfile to unix:/alloc/pulse/socket
echo "[svc-baresip] Using PULSE_SERVER=$PULSE_SERVER"

# Create baresip config directory
BARESIP_DIR="/root/.baresip"
mkdir -p "$BARESIP_DIR"

# Generate /etc/asound.conf to map ALSA devices to PulseAudio
# We use the 'type pulse' plugin (provided by libasound2-plugins)
# Append to existing config (created in Dockerfile)
cat <<EOF >> /etc/asound.conf
pcm.baresip_player {
    type pulse
    device "VirtualMic"
}

pcm.baresip_recorder {
    type pulse
    device "VirtualSpeaker.monitor"
}

ctl.baresip_player {
    type pulse
    device "VirtualMic"
}

ctl.baresip_recorder {
    type pulse
    device "VirtualSpeaker.monitor"
}
EOF

# Find module path — search for alsa.so instead of pulse.so
MOD_FILE=$(find /usr/lib -name alsa.so -path "*/baresip/*" 2>/dev/null | head -n 1)
if [ -z "$MOD_FILE" ]; then
    # Fallback: search for any baresip modules directory
    MOD_FILE=$(find /usr/lib -name "g711.so" -path "*/baresip/*" 2>/dev/null | head -n 1)
fi
if [ -z "$MOD_FILE" ]; then
    echo "[svc-baresip] ERROR: Could not find any baresip modules"
    echo "[svc-baresip] Debug: searching for anything baresip-related..."
    find /usr/lib -path "*baresip*" 2>/dev/null || echo "  nothing found under /usr/lib"
    dpkg -L baresip 2>/dev/null | head -20 || echo "  dpkg -L baresip failed"
    dpkg -L baresip-core 2>/dev/null | head -20 || echo "  dpkg -L baresip-core failed"
    sleep 30
    exit 1
fi
MOD_PATH=$(dirname "$MOD_FILE")
export MOD_PATH
echo "[svc-baresip] Module path: $MOD_PATH"
echo "[svc-baresip] Available modules: $(ls "$MOD_PATH"/*.so 2>/dev/null | xargs -n1 basename)"

# Generate config from template using envsubst
envsubst < /etc/baresip/config.template > "$BARESIP_DIR/config"
envsubst < /etc/baresip/accounts.template > "$BARESIP_DIR/accounts"

echo "[svc-baresip] Generated config:"
cat "$BARESIP_DIR/config"
echo "[svc-baresip] Starting Baresip..."

# Run baresip and process output for call detection
# Output goes to stdout (visible in docker logs) and we monitor for call events
# Use stdbuf to force line buffering so output appears immediately
stdbuf -oL baresip -f "$BARESIP_DIR" -v 2>&1 | while IFS= read -r line; do
    # Print line to stdout (docker logs will capture this)
    printf '%s\n' "$line"

    # Check for call events
    case "$line" in
        *"call: answering call"*)
            echo "[svc-baresip] DETECTED CALL! Creating trigger..."
            touch /alloc/call_active
            chmod 777 /alloc/call_active
            ;;
        *"terminated"*)
            echo "[svc-baresip] Call ended. Removing trigger..."
            rm -f /alloc/call_active
            ;;
    esac
done

# If we get here, baresip exited - the script will exit and s6 will restart it
echo "[svc-baresip] Baresip exited, will be restarted by s6"
