# rootfs/etc/s6-overlay/s6-rc.d/svc-baresip/run
#!/command/with-contenv bash
# shellcheck shell=bash
set -e

# PULSE_SERVER is set in Dockerfile to unix:/alloc/pulse/socket
echo "[svc-baresip] Using PULSE_SERVER=$PULSE_SERVER"

# Create baresip config directory
BARESIP_DIR="/root/.baresip"
mkdir -p "$BARESIP_DIR"

# Find module path
MOD_FILE=$(find /usr/lib -name pulse.so 2>/dev/null | grep baresip | head -n 1)
if [ -z "$MOD_FILE" ]; then
    echo "[svc-baresip] ERROR: Could not find baresip pulse.so module"
    exit 1
fi
MOD_PATH=$(dirname "$MOD_FILE")
export MOD_PATH

# Generate config from template using envsubst
envsubst < /etc/baresip/config.template > "$BARESIP_DIR/config"
envsubst < /etc/baresip/accounts.template > "$BARESIP_DIR/accounts"

echo "[svc-baresip] Generated config:"
cat "$BARESIP_DIR/config"
echo "[svc-baresip] Starting Baresip..."

# Run baresip and process output for call detection
# Output goes to stdout (visible in docker logs) and we monitor for call events
# Use stdbuf to force line buffering so output appears immediately
stdbuf -oL baresip -f "$BARESIP_DIR" -v 2>&1 | while IFS= read -r line; do
    # Print line to stdout (docker logs will capture this)
    printf '%s\n' "$line"

    # Check for call events
    case "$line" in
        *"call: answering call"*)
            echo "[svc-baresip] DETECTED CALL! Creating trigger..."
            touch /alloc/call_active
            chmod 777 /alloc/call_active
            ;;
        *"terminated"*)
            echo "[svc-baresip] Call ended. Removing trigger..."
            rm -f /alloc/call_active
            ;;
    esac
done

# If we get here, baresip exited - the script will exit and s6 will restart it
echo "[svc-baresip] Baresip exited, will be restarted by s6"
